# crontab
*/2 * * * * root /trigger_command.sh >> /var/log/cron.log 2>&1